# root/etc/services.d/forked/run
#!/usr/bin/with-contenv bash

echo 'Exec PlexVera...'

exec \
  php /plex-vera-app/bin/plexVera.php
